# cron.hourly/tattle-tale-status.sh
#!/bin/sh

datestring=$(date)
echo "TattleTale status at $datestring: RUNNING"
